# sat-api/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20

[program:sat-runner]
command=php /var/www/artisan sat:runner --loop
autostart=true
autorestart=true
priority=30
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/runner.log
